# openpower/custom/overlay/simics/etc/network/if-up.d/ifcfg-eth0
BOOTPROTO="dhcp"
STARTMODE="auto"
NAME="eth0"
